# porteus/modules/sheepit-module-xfce/home/guest/.config/autostart/TaskManager.desktop
[Desktop Entry]
Exec=xfce4-taskmanager
Name=TaskManager
StartupNotify=true
Type=Application
Version=1.0
